# contrib/gribble-supervisor.conf
[program:beaver]
directory=/tmp
command=/usr/local/bin/beaver
process_name=%(program_name)s%(process_num)02d
numprocs=1
autostart=true
autorestart=unexpected
startsecs=2
startretries=3
exitcodes=0,2
stopsignal=QUIT
stopwaitsecs=10
user=root
redirect_stderr=true
stdout_logfile=/var/log/beaver.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
environment=PYTHONUNBUFFERED=1,BEAVER_FILES="/var/log/syslog",BEAVER_TRANSPORT="redis",REDIS_NAMESPACE='logstash',REDIS_URL="redis://redis:6379/0"
